# supervisor.conf
[program:locustio-app]
command = locust -f /test/app/source/locustfile.py --host=%(ENV_HOST)s --logfile=/tmp/locusio.log

[program:nginx-app]
command = /usr/sbin/nginx
